# scripts/99-scan-plugged-usb.rules
# Action on the mounted partition
ACTION=="add", KERNEL=="sd*", RUN+="/path_to_your_script/scan-usb.sh"